# sv/slave/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

exec 2>&1
exec /spark*/bin/spark-class org.apache.spark.deploy.worker.Worker ${SPARK_MASTER_URL='spark://spark-master:7077'}
